# ddnss.timer
[Unit]
Description=Hourly update ddnss.de entries

[Timer]
OnCalendar=hourly
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
